# src/daemon/systemd/user/pulseaudio.service.in
[Unit]
Description=Sound Service
Requires=pulseaudio.socket
After=pulseaudio.socket

[Service]
# Note that notify will only work if --daemonize=no
Type=notify
ExecStart=@PA_BINARY@ --daemonize=no
Restart=on-failure

[Install]
Also=pulseaudio.socket
WantedBy=default.target
